# packaging/macos/com.copypaste.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
    LEGACY / OPT-IN ONLY (ADR-014, app-owned daemon lifecycle, 2026-05-30).

    The default install does NOT use this LaunchAgent. The CopyPaste desktop app
    starts the daemon as a child process on launch and SIGTERMs it on quit, and
    proactively `launchctl bootout`s this label on both launch and quit so a
    leftover agent can never resurrect the daemon behind the app's back.

    This plist is retained only for power users who want a headless,
    CLI-managed daemon WITHOUT the desktop app, via
    scripts/launchd/install-agent.sh (which kickstarts it explicitly). Do not
    load it alongside the app — the app will boot it out.
-->
<dict>
    <key>Label</key>
    <string>com.copypaste.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/CopyPaste.app/Contents/MacOS/copypaste-daemon</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!--
        Respawn only on crash; a clean exit (user invoked `daemon stop` or
        the daemon shut down intentionally) must NOT trigger a relaunch.
        Without this, `KeepAlive=<true/>` would fight `launchctl bootout`
        and the UI's quit/restart flows.
    -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>@HOME@/Library/Logs/CopyPaste/daemon.out.log</string>

    <key>StandardErrorPath</key>
    <string>@HOME@/Library/Logs/CopyPaste/daemon.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>

    <!--
        Bumped from 10 → 30 so a flapping daemon (e.g. config error,
        permission denial) doesn't burn CPU in a tight respawn loop.
    -->
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
